# mask-api.service.example
[Unit]
Description=Mask API
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/mask-api
ExecStart=/usr/bin/python3 /opt/mask-api/server.py
Restart=always
RestartSec=3
Environment=MASK_API_HOST=0.0.0.0
Environment=MASK_API_PORT=8080
Environment=MASK_PUBLIC_BASE_URL=http://111.228.12.141:8080

[Install]
WantedBy=multi-user.target
